# deploy/portal-messenger-prod.service
[Unit]
Description=Portal Znackar Messenger Worker (prod)
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/www/hosting/portalznackare.cz/www
Environment=APP_ENV=prod
ExecStart=/usr/bin/php bin/console messenger:consume async --time-limit=3600 --memory-limit=256M -v
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=portal-messenger-prod

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/www/hosting/portalznackare.cz/www/var

[Install]
WantedBy=multi-user.target